# roles/kubernetes/templates/systemd/kubelet.timer.j2
[Unit]
Description=Deley to start kubernetes

[Timer]
OnBootSec=5

[Install]
WantedBy=timers.target
